# supervisor.conf
; Luo OS — Supervisor Process Manager
; Start all services: supervisord -c supervisor.conf
; Control:           supervisorctl -c supervisor.conf status

[supervisord]
logfile=%(here)s/.luo_logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=%(here)s/.luo_logs/supervisord.pid
nodaemon=true
silent=false

[unix_http_server]
file=%(here)s/.luo_logs/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/.luo_logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:agent_api]
command=python3 %(here)s/ai_core/agent_api.py
directory=%(here)s
autostart=true
autorestart=true
startsecs=2
stderr_logfile=%(here)s/.luo_logs/agent_api.err.log
stdout_logfile=%(here)s/.luo_logs/agent_api.out.log
stdout_logfile_maxbytes=5MB

[program:rest_api]
command=python3 %(here)s/ai_core/rest_api.py
directory=%(here)s
autostart=true
autorestart=true
startsecs=2
stderr_logfile=%(here)s/.luo_logs/rest_api.err.log
stdout_logfile=%(here)s/.luo_logs/rest_api.out.log
stdout_logfile_maxbytes=5MB

[program:luo_agent_daemon]
command=python3 %(here)s/luo_agent/luo_agent.py --daemon
directory=%(here)s/luo_agent
autostart=true
autorestart=true
startsecs=5
stderr_logfile=%(here)s/.luo_logs/luo_agent.err.log
stdout_logfile=%(here)s/.luo_logs/luo_agent.out.log
stdout_logfile_maxbytes=5MB

[group:luo_os]
programs=agent_api,rest_api,luo_agent_daemon
