# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:qemu]
command=qemu-system-i386
    -cdrom /luo_os/luo_os.iso
    -serial tcp:0.0.0.0:7777,server,nowait
    -display vnc=:0
    -vga std
    -m 128M
    -no-reboot
    -no-shutdown
autorestart=true
stdout_logfile=/var/log/qemu.log
stderr_logfile=/var/log/qemu.log

[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log

[program:ollama]
command=ollama serve
autorestart=true
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama.log

[program:agent]
command=python3 /luo_os/agent/daemon.py
autorestart=true
stdout_logfile=/var/log/agent.log
stderr_logfile=/var/log/agent.log
